# files/etc/systemd/system/vncp-nginx-generator.service
[Unit]
Description=Generate nginx reverse-proxy config from VersaNode container labels
Wants=docker.service
After=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/vncp-nginx-generator.sh
