# files/etc/systemd/system/vncp-hostname.path
[Unit]
Description=Watch /etc/hostname and trigger VNCP TLS reissue

[Path]
PathChanged=/etc/hostname

[Install]
WantedBy=multi-user.target
