# 02-versanode-nginx-install/files/vncp-nginx-generate.service
[Unit]
Description=Generate Nginx config for VNCP proxies and Cockpit (HTTPS)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
EnvironmentFile=-/etc/default/vncp-nginx
ExecStart=/usr/local/sbin/vncp-nginx-generate
